# etc/init.d/debugfs
#!/sbin/openrc-run

description="Mounts the debugfs virtual filesystem"

depend() {
	# debugfs is mounted inside /sys, so it needs the sysfs service to run first.
	need sysfs
}

start() {
	# The mountinfo command is an OpenRC helper to check if a path is mounted.
	if mountinfo -q /sys/kernel/debug; then
		return 0
	fi

	ebegin "Mounting debugfs"
	mount -t debugfs none /sys/kernel/debug
	eend $?
}

stop() {
	# Unmounting is handled globally by the 'umountfs' shutdown service.
	# A manual stop action is not needed here.
	return 0
}