# etc/init.d/debugfs
#!/sbin/openrc-run

description="Mounts the debugfs virtual filesystem"

depend() {
	need sysfs
}

start() {
	if mountinfo -q /sys/kernel/debug; then
		return 0
	fi

	ebegin "Mounting debugfs"
	mount -t debugfs debugfs /sys/kernel/debug
	eend $?
}
