# debian/postinst
#!/usr/bin/env bash

# -- Add service to runlevel.

rc-update add debugfs-mount-openrc default

# -- Update services list.

rc-update -u
